# deploy/systemd/loop-duel-canary.timer
[Unit]
Description=Run LOOP DUEL canary hourly

[Timer]
OnBootSec=10m
OnUnitActiveSec=1h
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
