# deploy/systemd/loop-duel-canary.service
[Unit]
Description=LOOP live two-wallet DUEL canary
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/loop/current
EnvironmentFile=/opt/loop/shared/.env.production
Environment=PATH=/opt/loop/tools:/root/.acton/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/usr/bin/python3 scripts/run-duel-canary.py --contract ${LOOP_DUEL_CONTRACT_ADDRESS} --first-wallet loop-canary-a --second-wallet loop-canary-b
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/root/.acton /opt/loop/current/build

[Install]
WantedBy=multi-user.target
