# deploy/systemd/loop-duel-monitor.timer
[Unit]
Description=Check LOOP DUEL invariants every two minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=2m
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
